# tools/crai-update-database.timer
[Unit]
Description=CRAI update-database

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
